# schedule.cron
*/15 9-16 * * 1-5 python3 /home/<folder>/main.py